# infrastructure/system/systemd/arachne-scheduled-scrapes.timer
[Unit]
Description=Run Arachne scheduled scrapes

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=2m

[Install]
WantedBy=timers.target
